# supervisor/big_data.ini
[program:big_data]
command=/home/anaconda3/bin/gunicorn  app:app -c /usr/local/src/big_data/supervisor/gunicorn_config.py
stopasgroup=true
killasgroup=true
directory=/usr/local/src/big_data/
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/home/log_sellbot/big_data/supervisord-10001.log